# crontab
 * * * * * /usr/local/bin/python /coderbyte_ishango/src/scraping/results.py > /proc/1/fd/1 2>/proc/1/fd/2
